# conf/s6/services/php-fpm/run
#!/usr/bin/with-contenv sh
set -e;

/usr/sbin/php-fpm${PHP_VERSION} --nodaemonize --fpm-config /etc/php/7.3/fpm/php-fpm.conf