# conf/s6/services/nginx/run
#!/usr/bin/with-contenv bash
set -e;

if [ "$CONTAINERROLE" == "web" ]; then
     /usr/sbin/nginx
fi
