# conf/s6/services/nginx/run
#!/usr/bin/with-contenv bash
set -e;

if [ "$WEBSRV" == "nginx" ]; then
    /usr/sbin/nginx
fi
