# conf/s6/services/cronlog/run
#!/usr/bin/with-contenv bash

tail -f /var/log/cronlog