# conf/s6/services/cronlog/run
#!/usr/bin/with-contenv bash

if [ "$CONTAINERROLE" == "worker" ]; then
  tail -qf --follow=name --retry /var/log/cronlog
fi